# app/conf/supervise-celery.conf
[program:celery]
command=/bin/bash -c "cd /app && celery -A consumer worker --concurrency=5 --loglevel=info"
autorestart=true
