# app/conf/supervise-producer.conf
[program:producer]
command=/bin/bash -c "python /app/producer.py"
autorestart=false
